# .backup/local.d/04trackpoint.start
echo 250 > /sys/devices/platform/i8042/serio1/speed
echo 250 > /sys/devices/platform/i8042/serio1/sensitivity
